# deploy/basha.service
[Unit]
Description=basha — emotionally unstable terminal server
After=network.target

[Service]
Type=simple
User=basha
WorkingDirectory=/opt/basha
EnvironmentFile=/opt/basha/.env
ExecStart=/opt/basha/.venv/bin/python3 server.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

